# roles/repo/files/certbot-renew.service
[Unit]
Description=certbot renew certs
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew -q
